# thermostat.service
[Unit]
Description=Raspberry Pi Thermostat Control Service
After=network.target

[Service]
# IMPORTANT: Change 'landon' to your actual Raspberry Pi username if different
User=landon
# IMPORTANT: Change this to the absolute path of your project directory
WorkingDirectory=/home/landon/Thermostat 

# Type of service: simple means the script runs directly
Type=simple

# The command to execute: uses start_thermostat.sh which automatically detects venv
# The wrapper script will use venv/bin/python3 if it exists, otherwise system python3
# IMPORTANT: Update this path to match your project directory
ExecStart=/home/landon/Thermostat/start_thermostat.sh

# Restart the service if it stops for any reason
Restart=always
# Wait for 5 seconds before attempting to restart
RestartSec=5

[Install]
# Start the service automatically when the system boots up
WantedBy=multi-user.target
